# templates/mtproxy-panel-api.service
[Unit]
Description=MTProxy Logs Panel API
After=network.target mariadb.service
Wants=network.target

[Service]
Type=simple
User=mtpanel
Group=mtpanel
WorkingDirectory=__API_DIR__
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

